# btrfs-snapshot.service
[Unit]
Description=BTRFS Snapshots

[Service]
Type=oneshot
ExecStart=/usr/bin/btrfs-snapshot-manager schedule run
